# ikar-organizer.service
[Unit]
Description=IKAR RFID Organizer backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/alex/organizer_app
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/alex/organizer_app/service_backend.sh
ExecStopPost=/usr/bin/python3 /home/alex/organizer_app/cleanup_hardware.py
Restart=always
RestartSec=3
StartLimitIntervalSec=0
KillSignal=SIGTERM
TimeoutStopSec=12

[Install]
WantedBy=multi-user.target
