# packaging/deb/prerm
#!/bin/sh
set -e

SERVICE_NAME="rfid-local-mvp.service"

if command -v systemctl >/dev/null 2>&1; then
  systemctl stop "$SERVICE_NAME" || true
  systemctl disable "$SERVICE_NAME" >/dev/null 2>&1 || true
fi
