# packaging/deb/rfid-local-mvp.service
[Unit]
Description=RFID Local MVP backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/rfid-local-mvp
ExecStart=/usr/bin/rfid-local-mvp-backend
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
